# scripts/systemd/zkm-whatsapp-decrypt.service
[Unit]
Description=Decrypt WhatsApp msgstore backup into the zkm inbox, then zkm convert
Documentation=https://github.com/zommuter/zkm-whatsapp
# A Syncthing sync delivers the daily backup + several increment/tmp files in a
# burst, so the .path watcher fires this service many times within a second. Each
# fire is a cheap no-op (the wrapper's flock + newer-than guard), but systemd's
# default start-rate-limiter (5 starts / 10s) would trip and fail the .path unit,
# silently stopping the watcher. Disable the limiter — our own guards are the
# idempotency mechanism. Verified 2026-06-16 (W10 d058 live journey).
StartLimitIntervalSec=0

[Service]
Type=oneshot
# Copy zkm-whatsapp-decrypt.sh to ~/.local/bin (or point ExecStart at this repo copy).
# Configure paths/keyring via WA_* env here or in the script's defaults.
ExecStart=%h/.local/bin/zkm-whatsapp-decrypt.sh
TimeoutStartSec=10min
# Example overrides (uncomment + edit, or set defaults in the script):
# Environment=WA_CRYPT15=%h/Sync/WhatsApp/Databases/msgstore.db.crypt15
# Environment=WA_KEYREF=keyring:zkm/whatsapp:backup
# Environment=ZKM_BIN=%h/.local/bin/zkm
